# ec2/crontab
*/1 * * * * ~/lichess-swiss-league-website/lichessswissleague.sh start