# ec2/crontab
*/1 * * * * ~/lichess-swiss-league-website/script.sh start